# grapheqd.openrc
#!/sbin/openrc-run

: ${GRAPHEQD_BINARY:=${RC_PREFIX%/}/usr/local/sbin/grapheqd}
: ${GRAPHEQD_PIDFILE:=${RC_PREFIX%/}/run/${RC_SVCNAME}.pid}
: ${GRAPHEQD_ADDRESS:=0.0.0.0}
: ${GRAPHEQD_PORT:=8083}
: ${GRAPHEQD_SOUNDCARD:=hw,0}
: ${GRAPHEQD_USER:=nobody}

command="${GRAPHEQD_BINARY}"
pidfile="${GRAPHEQD_PIDFILE}"
command_args="-a ${GRAPHEQD_ADDRESS} -l ${GRAPHEQD_PORT} -s ${GRAPHEQD_SOUNDCARD} -u ${GRAPHEQD_USER} -p ${pidfile}"

depend() {
	need localmount
	after alsasound
}
